# systemd/awg-uplink-geo-domain-refresh.service
[Unit]
Description=Refresh geo domain lists and prepare nftset config
After=network-online.target awg-webui-ifaces.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=AWG_WEBUI_CFG_DIR=/etc/awg-uplink-webui
Environment=AWG_WEBUI_IFACE_ENV=/etc/awg-uplink-webui/interfaces.env
ExecStart=/usr/local/sbin/awg-uplink-geo-domain-refresh.py
