# systemd/awg-uplink-geo-routing-refresh.timer
[Unit]
Description=Periodic refresh of geo prefix list (allyouneed)

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=awg-uplink-geo-routing-refresh.service

[Install]
WantedBy=timers.target

