# systemd/awg-uplink-geo-routing-refresh.service
[Unit]
Description=Refresh geo prefix list and reapply selective routing
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# 1) Быстрый apply из локального кэша (если есть) — ошибки не фатальны.
ExecStart=-/usr/local/sbin/awg-uplink-geo-routing.sh apply-list
# 2) Фоновое обновление списка из сети.
ExecStart=/usr/local/sbin/awg-uplink-geo-routing-fetch.sh
# 3) Применение после обновления.
ExecStart=/usr/local/sbin/awg-uplink-geo-routing.sh apply-list

