# systemd/awg-uplink-geo-domain-nft-rotate.service
[Unit]
Description=Weekly geo domain nft set rotate (backup snapshot, flush primary)
After=network-online.target

[Service]
Type=oneshot
Environment=AWG_WEBUI_CFG_DIR=/etc/awg-uplink-webui
ExecStart=/usr/local/sbin/awg-uplink-geo-domain-refresh.py --rotate-nft
