# systemd/awg-uplink-geo-ipset-rotate.timer
[Unit]
Description=Daily awg-uplink geo ipset rotation (staging → route, flush staging)

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
