# systemd/awg-uplink-firewall.service
[Unit]
Description=Apply interface-panel firewall (nft awg_webui_fw or UFW rules)
After=network-online.target awg-webui-ifaces.service
Wants=network-online.target

[Service]
Type=oneshot
Environment=AWG_WEBUI_CFG_DIR=/etc/awg-uplink-webui
ExecStart=/usr/local/sbin/awg-uplink-firewall-apply.py
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
