# systemd/awg-quick@awg-uplink.service.d/override.conf
# awg-uplink после Docker; перед up ждём контейнер Amnezia (имя по умолчанию ^amnezia-awg — см. AMNEZIA_DOCKER_NAME_PATTERN в awg-uplink-policy.env).
# Ставится lib/awg-inject-uplink-policy.sh из репозитория

[Unit]
After=docker.service
Wants=docker.service

[Service]
# Не валим unit, если docker нет или контейнер переименован — hook всё равно ждёт порт в NAT.
ExecStartPre=-/bin/bash -c 'command -v docker >/dev/null 2>&1||exit 0; for i in $(seq 1 45); do docker ps --format "{{.Names}}" 2>/dev/null | grep -Eq "^amnezia-awg" && exit 0; sleep 2; done; exit 0'
